# software/src/10-mate.rules
KERNEL=="sd*", SUBSYSTEMS=="scsi",RUN+="/usr/local/bin/udev-handler.sh start"
