# software/src/10-mate.rules
KERNEL=="sd*", SUBSYSTEMS=="scsi",RUN+="/tmp/matemat"
